# root/etc/services.d/serverpackcreator/run
#!/usr/bin/with-contenv bash

cd /app/serverpackcreator

exec s6-setuidgid abc /usr/bin/java -jar serverpackcreator.jar -dock
