# root/etc/services.d/serverpackcreator/run
#!/usr/bin/with-contenv bash

cd /app/serverpackcreator

exec \
  s6-setuidgid abc \
    java -jar \
      /app/serverpackcreator/serverpackcreator.jar $START && \
  exit 0
